# deploy/kingshot-event-notification.service
[Unit]
Description=Kingshot Event Notification Discord Bot
After=network.target

[Service]
User=ubuntu
WorkingDirectory=/home/ubuntu/KingShot-Event-Notification
ExecStart=/home/ubuntu/KingShot-Event-Notification/venv/bin/python bot.py
EnvironmentFile=/home/ubuntu/KingShot-Event-Notification/.env
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
